# scripts/micro_ros_agent.service
[Unit]
Description=micro-ROS Agent (serial)
After=network.target

[Service]
Type=simple
User=pi
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/micro_ros_agent serial --dev /dev/pico_micro_ros -b 115200
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

# Copy this file to /etc/systemd/system/micro_ros_agent.service and run:
# sudo systemctl daemon-reload
# sudo systemctl enable --now micro_ros_agent.service
